# templates/gitolite.service.erb
[Unit]
Description=gitolite

[Service]
Type=oneshot
RemainAfterExit=yes
KillMode=none
User=git
ExecStart=cd /home/git/repositories/gitolite-admin.git
ExecStart=GL_LIBDIR=$(gitolite query-rc GL_LIBDIR) hooks/post-update refs/heads/master
ExecStop=/bin/true
