# templates/gitolite.service.erb
[Unit]
Description=gitolite
Requires=local-fs.target

[Service]
Type=oneshot
RemainAfterExit=yes
KillMode=none
User=git
WorkingDirectory=/home/git/repositories/gitolite-admin.git
Environment="GL_LIBDIR=/usr/share/gitolite3/lib"
ExecStart=/home/git/repositories/gitolite-admin.git/hooks/post-update refs/heads/master
ExecStop=/bin/true
